# ref/validator.conf
[program:validator]
command = sh /srv/validator-api/config/gunicorn_start.sh
user = ubuntu
stdout_logfile = /var/log/validator/gunicorn_supervisor.log
redirect_stderr = true
